# docker/base/rootfs/etc/s6-overlay/s6-rc.d/init-paseo/run
#!/command/with-contenv bash
# Seed PASEO_HOME and align the paseo user with the host volume owner.
set -euo pipefail

PUID="${PUID:-911}"
PGID="${PGID:-911}"

if [[ "$(id -g paseo)" != "$PGID" ]]; then
  groupmod -o -g "$PGID" paseo || true
fi
if [[ "$(id -u paseo)" != "$PUID" ]]; then
  usermod -o -u "$PUID" paseo || true
fi

SUDOERS_FILE="/etc/sudoers.d/paseo-enable-sudo"
case "$(printf '%s' "${PASEO_ENABLE_SUDO:-}" | tr '[:upper:]' '[:lower:]')" in
  true)
    mkdir -p /etc/sudoers.d
    printf 'paseo ALL=(ALL) NOPASSWD:ALL\n' > "$SUDOERS_FILE"
    chmod 0440 "$SUDOERS_FILE"
    echo "[init] PASEO_ENABLE_SUDO enabled: passwordless sudo granted to paseo."
    ;;
  *)
    rm -f "$SUDOERS_FILE"
    echo "[init] PASEO_ENABLE_SUDO disabled: passwordless sudo not granted."
    ;;
esac

IMAGE_DEFAULT_HOME="/home/paseo"
IMAGE_DEFAULT_PASEO_HOME="${IMAGE_DEFAULT_HOME}/.paseo"
CONTAINER_HOME="${HOME:-$IMAGE_DEFAULT_HOME}"
if [[ -z "${PASEO_HOME:-}" || ( "$CONTAINER_HOME" != "$IMAGE_DEFAULT_HOME" && "${PASEO_HOME:-}" == "$IMAGE_DEFAULT_PASEO_HOME" ) ]]; then
  PASEO_HOME="${CONTAINER_HOME}/.paseo"
fi
usermod -d "$CONTAINER_HOME" paseo || true

home_scoped_dir() {
  local name="$1"
  local suffix="$2"
  local image_default="${IMAGE_DEFAULT_HOME}${suffix}"
  local value="${!name:-}"
  if [[ -z "$value" || ( "$CONTAINER_HOME" != "$IMAGE_DEFAULT_HOME" && "$value" == "$image_default" ) ]]; then
    printf -v "$name" '%s' "${CONTAINER_HOME}${suffix}"
  fi
}

home_scoped_dir CLAUDE_CONFIG_DIR "/.claude"
home_scoped_dir CODEX_HOME "/.codex"
home_scoped_dir XDG_CONFIG_HOME "/.config"
home_scoped_dir XDG_DATA_HOME "/.local/share"
home_scoped_dir XDG_STATE_HOME "/.local/state"
home_scoped_dir XDG_CACHE_HOME "/.cache"

ensure_owned_dir() {
  local dir="$1"
  local mode="${2:-700}"
  if [[ -z "$dir" ]]; then
    return 0
  fi
  mkdir -p "$dir"
  chown paseo:paseo "$dir" || true
  chmod "$mode" "$dir" || true
}

# Only chown directory entries, not their contents; avoid clobbering large
# bind-mounted homes while still making the expected roots writable.
ensure_owned_dir "$CONTAINER_HOME" 700
ensure_owned_dir "$PASEO_HOME" 700
ensure_owned_dir "$CLAUDE_CONFIG_DIR" 700
ensure_owned_dir "$CODEX_HOME" 700
ensure_owned_dir "$XDG_CONFIG_HOME" 700
ensure_owned_dir "$XDG_DATA_HOME" 700
ensure_owned_dir "$XDG_STATE_HOME" 700
ensure_owned_dir "$XDG_CACHE_HOME" 700

echo "[init] Paseo paths ready (PUID=${PUID} PGID=${PGID} HOME=${CONTAINER_HOME} PASEO_HOME=${PASEO_HOME})"
