# docker/base/rootfs/etc/s6-overlay/s6-rc.d/init-paseo/run
#!/command/with-contenv bash
# Seed the /config volume and align the paseo user with the host volume owner.
set -euo pipefail

PUID="${PUID:-911}"
PGID="${PGID:-911}"

if [[ "$(id -g paseo)" != "$PGID" ]]; then
  groupmod -o -g "$PGID" paseo || true
fi
if [[ "$(id -u paseo)" != "$PUID" ]]; then
  usermod -o -u "$PUID" paseo || true
fi

SUDOERS_FILE="/etc/sudoers.d/paseo-enable-sudo"
case "$(printf '%s' "${PASEO_ENABLE_SUDO:-}" | tr '[:upper:]' '[:lower:]')" in
  true)
    mkdir -p /etc/sudoers.d
    printf 'paseo ALL=(ALL) NOPASSWD:ALL\n' > "$SUDOERS_FILE"
    chmod 0440 "$SUDOERS_FILE"
    echo "[init] PASEO_ENABLE_SUDO enabled: passwordless sudo granted to paseo."
    ;;
  *)
    rm -f "$SUDOERS_FILE"
    echo "[init] PASEO_ENABLE_SUDO disabled: passwordless sudo not granted."
    ;;
esac

PASEO_HOME="${PASEO_HOME:-/config}"
mkdir -p "$PASEO_HOME"

ensure_paseo_dir() {
  local dir="$1"
  local mode="${2:-700}"
  local existed=0
  if [[ -z "$dir" ]]; then
    return 0
  fi
  [[ -d "$dir" ]] && existed=1
  mkdir -p "$dir"
  case "$dir" in
    "$PASEO_HOME" | "$PASEO_HOME"/*)
      chown paseo:paseo "$dir" || true
      chmod "$mode" "$dir" || true
      ;;
    *)
      if [[ "$existed" == 0 ]]; then
        chown paseo:paseo "$dir" || true
        chmod "$mode" "$dir" || true
      fi
      ;;
  esac
}

# Only chown the top level + small paseo-owned config dirs; avoid clobbering a
# large bind-mounted /config on every boot.
ensure_paseo_dir "$PASEO_HOME" 700
ensure_paseo_dir "${CLAUDE_CONFIG_DIR:-${PASEO_HOME}/.claude}" 700
ensure_paseo_dir "${CODEX_HOME:-${PASEO_HOME}/.codex}" 700
ensure_paseo_dir "${XDG_CONFIG_HOME:-${PASEO_HOME}/.config}" 700
ensure_paseo_dir "${XDG_DATA_HOME:-${PASEO_HOME}/.local/share}" 700
ensure_paseo_dir "${XDG_CACHE_HOME:-${PASEO_HOME}/.cache}" 700

echo "[init] /config ready (PUID=${PUID} PGID=${PGID} home=${PASEO_HOME})"
