# docker/base/rootfs/etc/s6-overlay/s6-rc.d/svc-paseo/run
#!/command/with-contenv bash
# Run the daemon as the unprivileged paseo user.
exec s6-setuidgid paseo /usr/local/bin/paseo-run
