# docker/base/rootfs/etc/s6-overlay/s6-rc.d/init-docker-mods/run
#!/command/with-contenv bash
# Apply agent Docker Mods before the daemon starts so providers are detected.
exec /usr/local/bin/docker-mods
